# espruinohub/rootfs/etc/services.d/espruinohub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Hass.io Add-ons: Espruino Hub
# Starts Espruino Hub
# ==============================================================================
export BLENO_ADVERTISING_INTERVAL=300
export NOBLE_MULTI_ROLE=1
declare -a options

cd /opt/espruinohub || bashio::exit.nok "Could not change directory to Espruino Hub"
echo options > config.json

bashio::log.info "Starting Espruino Hub..."

exec node index.js